# deb/prerm
#!/bin/bash

# stop and remove service
systemctl stop advent-calendar.service
systemctl disable advent-calendar.service
